# media-sound/oss/files/oss
#!/sbin/runscript

depend() {
	need bootmisc localmount
	after modules isapnp coldplug hotplug
}

start() {
	ebegin "Starting OSS"
	/usr/sbin/soundon
	/usr/sbin/savemixer -L
	eend $? "Failed to start OSS"
}

terminate() {
	ebegin "Killing processes using OSS"
	devs=`perl -e '{while (<>) {m/^(\S*)/; print "/dev/$1\n"}}' \
		< /proc/opensound/devfiles`
	fuser -k ${devs}  >/dev/null 2>/dev/null
	eend 0
}

stop() {
	terminate
	ebegin "Stopping OSS"
	/usr/sbin/savemixer
	/usr/sbin/soundoff
	eend $? "Failed to stop OSS"
}
